# data/storaged.service.in
[Unit]
Description=LVM Disk Manager
Documentation=man:@storagedexecname@(8)

[Service]
Type=dbus
BusName=com.redhat.storaged
ExecStart=@storagedprivdir@/@storagedexecname@
